# autostart/steam-start.desktop
[Desktop Entry]
Encoding=UTF-8
Name=steam
Comment=Start steam on system boot
Exec=/usr/bin/steam -tenfoot
Terminal=false
Type=Application
